# examples/supervisor.example.conf
; AI Ops Worker - Supervisor Configuration
; Copy to /etc/supervisor/conf.d/ai-ops-worker.conf

[program:ai-ops-worker]
command=/srv/ai-ops/venv/bin/python /srv/ai-ops/worker.py
directory=/srv/ai-ops
user=ubuntu
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/ai-ops/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/ai-ops/worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=
    PATH="/srv/ai-ops/venv/bin:/usr/local/bin:/usr/bin:/bin",
    HOME="/home/ubuntu",
    CI="true",
    TERM="dumb"
